# Deployment/Systemd/qw-renew-letsencrypt.timer
[Unit]
Description=To renew Let's Encrypt daily

[Timer]
OnCalendar=*-*-* 4:16:00

[Install]
WantedBy=timers.target

